# cast_homeassistant.service
[Unit]
Description=Cast Home Assistant

[Service]
User=pi
Type=oneshot
ExecStart=/home/pi/homeassistant/cast.sh